# airootfs/etc/xdg/autostart/custom.desktop
# Start Hyprland in the live environment
[Desktop Entry]
Type=Application
Exec=Hyprland